# templates/exporter.service.j2
{{ ansible_managed | comment }}

[Unit]
Description=Prometheus {{ _exporter.name }}
Documentation=https://github.com/{{ _exporter_github_org }}/{{ _exporter_github_repo }}
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ _exporter_system_user }}
Group={{ _exporter_system_group }}
{% for key, value in _exporter_env_vars.items() %}
Environment="{{ key }}={{ value }}"
{% endfor %}
ExecStart={{ prometheus_exporter_install_dir }}/{{ _exporter_binary_name }} \
{% if _exporter.web_listen_flag is defined %}
    {{ _exporter.web_listen_flag }}={{ _exporter_listen_address }}:{{ _exporter_listen_port }} \
{% else %}
    --web.listen-address={{ _exporter_listen_address }}:{{ _exporter_listen_port }} \
{% endif %}
{% for arg in _exporter_extra_args %}
    {{ arg }}{% if not loop.last %} \{% endif %}

{% endfor %}

ExecReload=/bin/kill -HUP $MAINPID
SyslogIdentifier={{ _exporter.name }}
Restart=on-failure
RestartSec=5

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectHostname=yes
ProtectClock=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
ReadWritePaths={{ _exporter_config_dir }}

[Install]
WantedBy=multi-user.target
